# systemd/python-album-art-watcher.path
[Path]
PathModified=/opt/album-art/album-art.py
#PathChanged=/opt/album-art/album-art.py

[Install]
WantedBy=multi-user.target
