# systemd/python-album-art-watcher.path
[Path]
PathModified=/opt/python-album-art/album-art.py

[Install]
WantedBy=multi-user.target
